# deploy/buddharauer-backup.timer
[Unit]
Description=Run the Buddharauer backup daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
